# config/router/rc.local
#!/bin/sh

zero="192.168.1.2"
ports="500 587 993 1701 1723 4500"

for port in ${ports}; do
    iptables -t nat -I PREROUTING -p tcp -i eth1 \
        --dport ${port} -j DNAT --to ${zero}:${port}
    iptables -t nat -I PREROUTING -p udp -i eth1 \
        --dport ${port} -j DNAT --to ${zero}:${port}
done

exit 0
